# scripts/hwmond.service
[Unit]
Description=Apple Xserve Hardware Monitor (LED + BMC)
Documentation=https://github.com/mav2287/hwmond
After=network.target ipmi.service
Wants=network.target

[Service]
Type=simple
ExecStart=/usr/local/sbin/hwmond
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hwmond

# Security hardening
NoNewPrivileges=true
ProtectHome=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/dev /var/run
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictSUIDSGID=true
MemoryDenyWriteExecute=true
LockPersonality=true

[Install]
WantedBy=multi-user.target
